# blockchain-node.service
[Unit]
Description=Blockchain Rust P2P Node
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/tuyenpkt/blockchain-rust
ExecStart=/home/tuyenpkt/blockchain-rust/target/release/blockchain-rust node 8333
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=blockchain-node

# Giới hạn resource
LimitNOFILE=65535
MemoryMax=512M

[Install]
WantedBy=multi-user.target
