# pkt-sync.service
[Unit]
Description=PKT Testnet Header Sync
After=network.target blockchain-node.service
Wants=network-online.target

[Service]
Type=simple
User=tuyenpkt
WorkingDirectory=/home/tuyenpkt/blockchain-rust
# sync từ pkt-node local (8333), skip PoW check vì dùng chain riêng
ExecStart=/home/tuyenpkt/blockchain-rust/target/release/blockchain-rust sync --skip-pow
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pkt-sync

LimitNOFILE=65535
MemoryMax=256M

[Install]
WantedBy=multi-user.target
